# cron/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=rsyslogd -n

[program:cron]
command=cron -f -L 15
